# auditd-slave/templates/etc/init/audit.conf.tmpl
start on runlevel 2

exec /etc/init.d/auditd start

post-start script
        /sbin/auditctl -R /etc/audit/audit.rules
end script